# srv/foo.socket
[Unit]
Description=Socket activated by sd_notify

[Socket]
ListenDatagram=/run/custom_sd_notify_socket
