# srv/foo.service
[Service]
ExecStart=/usr/bin/echo "foo service is running"

