# srv/notify.service
[Unit]
Description=Service with type notify

[Service]
Restart=always
TimeoutStartSec=10
Type=notify
NotifyAccess=main
ExecStart=/opt/notify
