# srv/watchdog.service
[Unit]
Description=Watchdog Service

[Service]
Type=simple
Restart=always
ExecStart=/opt/watchdog
WatchdogSec=10
StartLimitBurst=0
